# config/imperium.service
[Unit]
Description=Imperium Intent-Based Networking Controller
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=imperium
Group=imperium
WorkingDirectory=/home/imperium/Imperium
Environment="PATH=/home/imperium/Imperium/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONPATH=/home/imperium/Imperium"
Environment="PYTHONUNBUFFERED=1"
Environment="FLASK_ENV=production"
Environment="NETWORK_INTERFACE=eth0"

# Ensure database and log directories exist
ExecStartPre=/bin/mkdir -p /home/imperium/Imperium/data
ExecStartPre=/bin/mkdir -p /home/imperium/Imperium/logs
ExecStartPre=/bin/chown -R imperium:imperium /home/imperium/Imperium/data
ExecStartPre=/bin/chown -R imperium:imperium /home/imperium/Imperium/logs

# Start the Imperium controller
ExecStart=/home/imperium/Imperium/venv/bin/python3 /home/imperium/Imperium/src/main.py

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Security settings
NoNewPrivileges=false
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=imperium

# Resource limits (updated for modern systemd)
MemoryMax=512M
CPUQuota=60%

[Install]
WantedBy=multi-user.target
